# pbs/autoshutdown/autoshutdown.service
[Unit]
Description=Shutdown when there are no tasks left

[Service]
ExecStart=/scripts/autoshutdown.sh

[Install]
WantedBy=default.target
