# lxc/failover_openhab/failover.service
[Unit]
Description=Autoshutdown

[Service]
Type=oneshot
ExecStart=/bin/bash /scripts/failover_openhab_master.sh
